# .cicd/webhook.service
[Unit]
Description=GitHub Webhook Listener for darth-forge
After=network.target

[Service]
Type=simple
User=darth
WorkingDirectory=/home/darth/darth-forge
Environment="WEBHOOK_PORT=9000"
Environment="DEPLOY_SCRIPT=/home/darth/darth-forge/.cicd/deploy.sh"
# WEBHOOK_SECRET is loaded from .cicd/.webhook_secret by the Python script
ExecStart=/bin/bash -c 'export WEBHOOK_SECRET=$(cat /home/darth/darth-forge/.cicd/.webhook_secret) && /usr/bin/python3 /home/darth/darth-forge/.cicd/webhook-listener.py'
Restart=always
RestartSec=10
StandardOutput=append:/home/darth/darth-forge/.cicd/webhook.log
StandardError=append:/home/darth/darth-forge/.cicd/webhook.log

[Install]
WantedBy=multi-user.target
